# lr9/mephi.path
[Unit]
Description="CHeck for MEPHI file in /tmp"
After=sshd.service
[Path]
PathExists=/tmp/mephi.file
Unit=mephi.service
[Install]
WantedBy=multi-user.target
