# lr9/mephi.service
[Unit]
Description=MEPhI hello service
After=sshd.service
[Service]
ExecStart=/usr/bin/printf "Hello, MephI"
[Install]
WantedBy=multi-user.target
